# xvfb-firefox/xvfb.sv.conf
[program:xvfb]
command=Xvfb :1 -screen 5 1024x768x8
